# cron/swif.cron
#MAILTO="bobcat@jlab.org"
*/30 * * * * /path/to/clas12-workflow/cron/swif-cron.csh
